# crontab
* * * * * /app/server php-cli artisan schedule:run >> /dev/null 2>&1
